# tools/gate/playbooks/build-helm-packages/templates/helm-serve.service.j2
[Unit]
Description=Helm Server
After=network.target

[Service]
User={{ helm_server_user.stdout }}
Restart=always
ExecStart=/usr/bin/helm serve

[Install]
WantedBy=multi-user.target
